# systemd/traefik-admin-helper.service
[Unit]
Description=Traefik Admin write helper
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=dan
Group=dan
WorkingDirectory=/home/dan/work/tmweb/infra/traefik
ExecStart=/usr/bin/python3 /home/dan/work/tmweb/infra/traefik/write-helper/traefik_admin_helper.py
Restart=on-failure
RestartSec=2
NoNewPrivileges=false
PrivateTmp=true
ProtectSystem=full
ProtectHome=false
ReadOnlyPaths=/home/dan/work/tmweb/infra/traefik/write-helper /home/dan/work/tools
ReadWritePaths=/srv/tmweb/traefik/conf /srv/tmweb/traefik/logs /srv/tmweb/traefik/run

[Install]
WantedBy=multi-user.target
